# deploy/systemd/cowork-os.service
[Unit]
Description=CoWork OS (Headless)
After=network.target

[Service]
Type=simple

# Path where you cloned/built CoWork OS on the VPS.
WorkingDirectory=/opt/ChatAndBuild

# Optional env file (LLM keys, data dir, control plane host/port, etc).
EnvironmentFile=-/etc/ChatAndBuild.env

# Headless Electron can fail on some servers without a DISPLAY. Xvfb provides a virtual display.
ExecStart=/usr/bin/xvfb-run -a /usr/bin/node /opt/ChatAndBuild/bin/ChatAndBuildd.js

Restart=on-failure
RestartSec=5

# Run as an unprivileged user (recommended).
User=ChatAndBuild
Group=ChatAndBuild

[Install]
WantedBy=multi-user.target
